# roles/hardening_os/templates/etc/profile.d/profile_hist_tmout.conf.j2
HISTSIZE=10000
TMOUT=300
HISTTIMEFORMAT="%d/%m/%y %T "
shopt -s histappend
